# files/supervisor-zerotier.conf
[supervisord]
[program:zerotier]
command=zerotier-one
redirect_stdout=false
redirect_stderr=false